# scanner/post-yum-security.cron
#!/bin/bash
/sbin/grubby --default-kernel | /bin/grep -q "$(/bin/uname -r)" || /sbin/reboot
