# scanner/post-yum-security.cron
#!/bin/bash
/usr/bin/needs-restarting -r > /dev/null || /usr/sbin/shutdown --reboot +5
