# scanner/vuls.cron
#!/bin/bash
cd /root/vuls
#./vuls.sh 123456789012 -lang=ja -to-slack > /dev/null
